# .config/systemd/user/ipfs.service
[Unit]
Description=IPFS Daemon (VPN-restricted via br0 bridge)
After=network-online.target
Wants=network-online.target
# Note: firejail-bridge.service is a system service that creates br0

[Service]
Type=simple
ExecStartPre=/bin/bash -c 'until [ -e /sys/class/net/br0 ]; do sleep 1; done'
ExecStart=/usr/bin/firejail --net=br0 --protocol=unix,inet,inet6,netlink --dns=10.64.0.1 /usr/bin/ipfs daemon
Restart=on-failure
RestartSec=10
Environment="IPFS_PATH=%h/.ipfs"

[Install]
WantedBy=default.target
