# .config/systemd/user/voxd-tray.service
[Unit]
Description=VOXD tray mode (user)
After=default.target

[Service]
Type=simple
WorkingDirectory=%h/repos/voxd
ExecStart=/usr/bin/uv run voxd --tray
Restart=on-failure
RestartSec=2s
Environment=YDOTOOL_SOCKET=%h/.ydotool_socket
Environment=QT_QPA_PLATFORM=wayland

[Install]
WantedBy=default.target
